# deploy/docker/php-fpm/supervisor.d/horizon.ini
[supervisord]
user=root
nodaemon=true

[program:horizon]
command=/usr/local/bin/php /app/artisan horizon
autorestart=true
priority=5
