# deploy/docker/php-fpm/supervisor.d/php-fpm.ini
[supervisord]
user=root
nodaemon=true

[program:php-fpm]
command=/usr/local/sbin/php-fpm -F
autorestart=true
priority=10