# deploy/docker/php-fpm/supervisor.d/schedule.ini
[supervisord]
user=root
nodaemon=true

[program:schedule]
command=sh -c "sleep 60 && /usr/local/bin/php /app/artisan schedule:run"
autorestart=true
priority=5